# docker/root/etc/s6-overlay/s6-rc.d/svc-paper-agent/run
#!/bin/bash
# Ensure /config is writable by abc so the DB can be created and .env checked
chown -R abc:abc /config

cd /app
# Check if .env exists in /config (the volume)
if [ ! -f /config/.env ]; then
    echo "No .env found in /config. Creating from example..."
    if [ -f .env.example ]; then
        cp .env.example /config/.env
    else
        echo "No .env.example found. Creating empty .env in /config"
        touch /config/.env
    fi
fi

# Fix potential wrong database path in .env (whether newly created or existing)
# If it points to /data/ (default), change it to /config/ (volume)
sed -i 's|sqlite:///data/paper_agent.db|sqlite:///config/paper_agent.db|g' /config/.env

# Fix permissions so user abc can read/write the .env file
chown abc:abc /config/.env

# Link /config/.env to /app/.env just in case, though config.py handles it
# ln -sf /config/.env /app/.env

export UV_PYTHON_INSTALL_DIR="/app/.uv_python"
export UV_CACHE_DIR="/app/.uv_cache"

# exec ensures the process replaces the shell, handling signals correctly
exec s6-setuidgid abc uv run uvicorn src.main:app --host 0.0.0.0 --port 8000
